# supervisord.conf
; supervisord runs as PID 1 (via tini) and manages code-server + hermes-agent
; + hermes-webui. All programs run as `coder`. stdout/stderr go to the pod
; log so `kubectl logs` shows everything.

[supervisord]
nodaemon=true
user=root
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/run/supervisord.pid

[unix_http_server]
file=/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[program:ttyd]
; Standalone xterm.js → bash bridge on :7681. Distinct from code-server's
; integrated terminal — gives users a lightweight "open a shell in the pod"
; URL without loading the full VSCode UI. Login shell so coder's profile
; (including PATH for /opt/hermes/.venv if they activate it) loads.
command=/usr/local/bin/ttyd --port 7681 --interface 0.0.0.0 --writable bash -l
user=coder
directory=/home/coder
environment=HOME="/home/coder",USER="coder",SHELL="/bin/bash",DOCKER_HOST="tcp://localhost:2375"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:code-server]
command=/usr/local/bin/code-server
    --bind-addr=0.0.0.0:8080
    --auth=none
    --disable-telemetry
    --disable-update-check
    /home/coder/workspace
user=coder
directory=/home/coder
environment=HOME="/home/coder",USER="coder",SHELL="/bin/bash",DOCKER_HOST="tcp://localhost:2375"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:hermes-agent]
command=/usr/local/bin/hermes-agent-bootstrap gateway run
user=coder
directory=/opt/hermes
environment=HOME="/home/coder",USER="coder",HERMES_HOME="/home/coder/.hermes",HERMES_WORKSPACE="/home/coder/workspace",DOCKER_HOST="tcp://localhost:2375",PATH="/opt/hermes/.venv/bin:/usr/local/bin:/usr/bin:/bin"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:hermes-webui]
; bootstrap.py is a one-shot launcher that fork-detaches server.py and exits —
; supervising it leads to FATAL because exit 0 ≠ "stay running". Run server.py
; directly so supervisord sees the actual long-lived process.
;
; Use the AGENT'S venv python (not /usr/bin/python3) — webui imports
; run_agent.AIAgent which transitively pulls in openai/anthropic/etc., all of
; which live in /opt/hermes/.venv/lib/python3.13/site-packages. Without this,
; webui boots but every chat request fails with "AIAgent not available".
command=/opt/hermes/.venv/bin/python3 /opt/hermes-webui/server.py
user=coder
directory=/opt/hermes-webui
environment=HOME="/home/coder",USER="coder",HERMES_WEBUI_HOST="0.0.0.0",HERMES_WEBUI_PORT="8787",HERMES_WEBUI_STATE_DIR="/home/coder/.hermes/webui",HERMES_HOME="/home/coder/.hermes",HERMES_WEBUI_AGENT_DIR="/opt/hermes",HERMES_WEBUI_DEFAULT_WORKSPACE="/home/coder/workspace"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
