# supervisord-desktop.conf
; Desktop-variant supervisord — same programs as the basic config plus
; KasmVNC (which itself starts Xvnc + the xfce4-session via xstartup).
;
; Architecture:
;   kasmvnc       owns DISPLAY=:1 (Xvnc + WebSocket on 6901). Reads
;                 ~/.vnc/xstartup which dbus-launches xfce4-session, so
;                 the desktop environment is up the moment the VNC display
;                 is ready.
;   hermes-agent  inherits DISPLAY=:1 from this config — headed Playwright
;                 windows + any xdotool/scrot calls render on the same X
;                 display the user sees through KasmVNC. Agent and human
;                 work on the same desktop.
;   ttyd          unchanged — still useful as a "raw shell" entrypoint.
;   code-server / hermes-webui — unchanged.

[supervisord]
nodaemon=true
user=root
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/run/supervisord.pid

[unix_http_server]
file=/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[program:kasmvnc]
; -fg keeps it foreground for supervisord; -interface 0.0.0.0 binds the
; WebSocket to all interfaces (the gateway is in front, no direct exposure).
; SecurityTypes None turns off VNC-level auth — the platform's oauth2-proxy
; already gates the gateway hostname. -BlacklistThreshold disables the
; default brute-force lockout (irrelevant when auth is offloaded).
command=/usr/bin/kasmvncserver
    -fg
    -interface 0.0.0.0
    -websocketPort 6901
    -SecurityTypes None
    -DisableBasicAuth
    -BlacklistThreshold 0
    :1
user=coder
directory=/home/coder
environment=HOME="/home/coder",USER="coder",SHELL="/bin/bash",DISPLAY=":1"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:ttyd]
; --client-option flags forwarded to xterm.js. See supervisord.conf
; for the basic image — kept in sync.
command=/usr/local/bin/ttyd --port 7681 --interface 0.0.0.0 --writable
    --client-option fontFamily="JetBrainsMono Nerd Font, Menlo, Consolas, monospace"
    --client-option fontSize=14
    --client-option lineHeight=1.2
    --client-option cursorBlink=true
    --client-option theme={"background":"#1e1e1e","foreground":"#d4d4d4","cursor":"#d4d4d4","selectionBackground":"#264f78","black":"#000000","red":"#cd3131","green":"#0dbc79","yellow":"#e5e510","blue":"#2472c8","magenta":"#bc3fbc","cyan":"#11a8cd","white":"#e5e5e5","brightBlack":"#666666","brightRed":"#f14c4c","brightGreen":"#23d18b","brightYellow":"#f5f543","brightBlue":"#3b8eea","brightMagenta":"#d670d6","brightCyan":"#29b8db","brightWhite":"#e5e5e5"}
    bash -l
user=coder
directory=/home/coder
environment=HOME="/home/coder",USER="coder",SHELL="/bin/bash",DOCKER_HOST="tcp://localhost:2375"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:code-server]
command=/usr/local/bin/code-server
    --bind-addr=0.0.0.0:8080
    --auth=none
    --disable-telemetry
    --disable-update-check
    /home/coder/workspace
user=coder
directory=/home/coder
environment=HOME="/home/coder",USER="coder",SHELL="/bin/bash",DOCKER_HOST="tcp://localhost:2375"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:hermes-agent]
; DISPLAY=:1 lets the agent open headed browsers and use xdotool/scrot on
; the same X server the user sees over KasmVNC. With KasmVNC owning the
; display, headed Playwright windows show up live.
;
; AGENT_BROWSER_HEADED=1 flips agent-browser (and its underlying Playwright)
; out of the default headless mode. Without this the agent's browser tool
; spawns a headless Chromium that never paints to the X server — user sees
; an empty desktop while tests "run". With it, every navigation/click is
; live on the KasmVNC view.
command=/usr/local/bin/hermes-agent-bootstrap gateway run
user=coder
directory=/opt/hermes
environment=HOME="/home/coder",USER="coder",HERMES_HOME="/home/coder/.hermes",HERMES_WORKSPACE="/home/coder/workspace",DOCKER_HOST="tcp://localhost:2375",DISPLAY=":1",AGENT_BROWSER_HEADED="1",PATH="/opt/hermes/.venv/bin:/usr/local/bin:/usr/bin:/bin"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:hermes-webui]
; bootstrap.py is a one-shot launcher that fork-detaches server.py and exits —
; supervising it leads to FATAL because exit 0 ≠ "stay running". Run server.py
; directly so supervisord sees the actual long-lived process.
;
; Use the AGENT'S venv python (not /usr/bin/python3) — webui imports
; run_agent.AIAgent which transitively pulls in openai/anthropic/etc., all of
; which live in /opt/hermes/.venv/lib/python3.13/site-packages. Without this,
; webui boots but every chat request fails with "AIAgent not available".
command=/opt/hermes/.venv/bin/python3 /opt/hermes-webui/server.py
user=coder
directory=/opt/hermes-webui
environment=HOME="/home/coder",USER="coder",HERMES_WEBUI_HOST="0.0.0.0",HERMES_WEBUI_PORT="8787",HERMES_WEBUI_STATE_DIR="/home/coder/.hermes/webui",HERMES_HOME="/home/coder/.hermes",HERMES_WEBUI_AGENT_DIR="/opt/hermes",HERMES_WEBUI_DEFAULT_WORKSPACE="/home/coder/workspace"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
